# anime_monitor.service
[Unit]
Description=Anime Monitor Telegram Bot
After=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/anime-monitor-bot
ExecStart=/root/anime-monitor-bot/venv/bin/python /root/anime-monitor-bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
